# proxy_pool/worker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:openvpn]
command=openvpn /var/run/secrets/config.ovpn
priority=100

[program:squid]
command=squid -N
priority=200

